# fs/etc/services.d/portainer-agent/run
#!/usr/bin/with-contenv ash

wait-for-signal || exit 10

# The executable should be running from the /app as defined in the original Dockerfile (WORKDIR directive is not honored by S6)
cd /app

echo "Starting  Portainer agent..."

exec s6-setuidgid "$PUID:$PGID" ./agent $PORTAINER_AGENT_ARGS