# fs/etc/services.d/portainer-agent/run
#!/usr/bin/with-contenv sh

s6-setuidgid nonroot

# Wait for signal file

wait-for-signal || exit 10

echo "Start Portainer agent..."

exec /app/agent $PORTAINER_AGENT_ARGS